# recipes-core/initscripts/files/90-rfkill.rules
SUBSYSTEM=="rfkill", ATTR{type}=="bluetooth", RUN+="/etc/init.d/udev_bt.sh"
