# recipes-multimedia/tvheadend/files/sv/tvheadend/run
#!/bin/sh

# NOTE: This presumes that you're adding the following user/group...
exec /usr/bin/tvheadend -C -d -u tvheadend -g tvheadend 2>&1
