# recipes-multimedia/tvheadend/files/sv/tvheadend/run
#!/bin/sh
exec /usr/bin/tvheadend 2>&1 
