# recipes-connectivity/mka-daemon/files/sv/mkad/run
#!/bin/sh

# Start the service with the pre-defined config file path and in foreground...
exec mkad -f -c /etc/mkad/mkad.conf

